# deploy/com.dizhang.cloudflared-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dizhang.cloudflared-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>mkdir -p "$HOME/.cloudflared/logs"; exec "$HOME/.cloudflared/cloudflared" --config "$HOME/.cloudflared/config.yml" tunnel run &gt;&gt; "$HOME/.cloudflared/logs/cloudflared.launchd.out.log" 2&gt;&gt; "$HOME/.cloudflared/logs/cloudflared.launchd.err.log"</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users</string>
</dict>
</plist>
