# deploy/com.zzddddzz.msds498-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zzddddzz.msds498-dashboard</string>
  <key>WorkingDirectory</key>
  <string>/Users</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>repo="$HOME/GitHub/insurance-outcome-intelligence-capstone"; mkdir -p "$repo/logs"; cd "$repo" &amp;&amp; exec "$repo/.venv/bin/python" -m streamlit run dashboard.py --server.port 8502 --server.address 127.0.0.1 --server.headless true --server.fileWatcherType none --browser.gatherUsageStats false &gt;&gt; "$repo/logs/msds498-dashboard.out.log" 2&gt;&gt; "$repo/logs/msds498-dashboard.err.log"</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
